# journal-triggerd.service.in
[Unit]
Description=Triggers on journal messages
After=systemd-journald.service

[Service]
ExecStart=@bindir@/journal-triggerd

[Install]
WantedBy=multi-user.target
